# raspberry/config/radioPowerOn.service
[Unit]
After=radioPowerOn.service

[Service]
ExecStart=/etc/zabbix/scripts/radioPowerOn.sh

[Install]
WantedBy=default.target
